# houzzkit-f1-bsp-k6199/yocto/meta-rockchip/recipes-bsp/bt-service/files/bt-attach.service

[Unit]
Description=bluetooth-toggle

[Service]
Type=forking
ExecStart=/usr/bin/bt-attach
