# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/lib/systemd/system/tpm-udev.service
[Unit]
Description=Handle dynamically added tpm devices

[Service]
ExecStart=systemd-tmpfiles --create tpm-udev.conf
