# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/lib/systemd/system/tpm-udev.path
[Unit]
ConditionVirtualization=container
Description=Handle dynamically added tpm devices

[Path]
PathChanged=/dev

[Install]
WantedBy=paths.target
