# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/share/gvfs/mounts/dns-sd.mount
[Mount]
Type=dns-sd
Exec=/usr/libexec/gvfsd-dnssd
AutoMount=true
DBusName=org.gtk.vfs.mountpoint_dnssd
