# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/share/gvfs/mounts/google.mount
[Mount]
Type=google-drive
Exec=/usr/libexec/gvfsd-google
AutoMount=false
