# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/share/gvfs/mounts/onedrive.mount
[Mount]
Type=onedrive
Exec=/usr/libexec/gvfsd-onedrive
AutoMount=false
