# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/share/gvfs/mounts/wsdd.mount
[Mount]
Type=wsdd
Exec=/usr/libexec/gvfsd-wsdd
AutoMount=true