# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/usr/lib/systemd/user/gvfs-metadata.service
[Unit]
Description=Virtual filesystem metadata service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/libexec/gvfsd-metadata
Type=dbus
BusName=org.gtk.vfs.Metadata
Slice=session.slice
