# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/avahi-daemon.list
/.
/etc
/etc/avahi
/etc/avahi/avahi-daemon.conf
/etc/avahi/hosts
/etc/avahi/services
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/avahi-daemon.service
/usr/lib/systemd/system/avahi-daemon.socket
/usr/sbin
/usr/sbin/avahi-daemon
/usr/share
/usr/share/avahi
/usr/share/avahi/avahi-service.dtd
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.AddressResolver.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.DomainBrowser.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.EntryGroup.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.HostNameResolver.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.RecordBrowser.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.Server.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceBrowser.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceResolver.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceTypeBrowser.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.Avahi.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/avahi-dbus.conf
/usr/share/doc
/usr/share/doc/avahi-daemon
/usr/share/doc/avahi-daemon/copyright
/usr/share/doc/avahi-daemon/examples
/usr/share/doc/avahi-daemon/examples/example.service
/usr/share/doc/avahi-daemon/examples/sftp-ssh.service
/usr/share/doc/avahi-daemon/examples/ssh.service
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/avahi-daemon.conf.5.gz
/usr/share/man/man5/avahi.hosts.5.gz
/usr/share/man/man5/avahi.service.5.gz
/usr/share/man/man8
/usr/share/man/man8/avahi-daemon.8.gz
/usr/share/doc/avahi-daemon/NEWS.gz
/usr/share/doc/avahi-daemon/README
/usr/share/doc/avahi-daemon/changelog.Debian.gz
