# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/bluez.conffiles
/etc/bluetooth/input.conf
/etc/bluetooth/main.conf
/etc/bluetooth/network.conf
/etc/default/bluetooth
/etc/init.d/bluetooth
