# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/chrony.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.chronyd
/etc/chrony
/etc/chrony/conf.d
/etc/chrony/conf.d/README
/etc/chrony/sources.d
/etc/chrony/sources.d/README
/etc/default
/etc/default/chrony
/etc/dhcp
/etc/dhcp/dhclient-exit-hooks.d
/etc/dhcp/dhclient-exit-hooks.d/chrony
/etc/init.d
/etc/init.d/chrony
/etc/logrotate.d
/etc/logrotate.d/chrony
/etc/network
/etc/network/if-post-down.d
/etc/network/if-post-down.d/chrony
/etc/network/if-up.d
/etc/network/if-up.d/chrony
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/chrony
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/chrony
/usr
/usr/bin
/usr/bin/chronyc
/usr/lib
/usr/lib/NetworkManager
/usr/lib/NetworkManager/dispatcher.d
/usr/lib/NetworkManager/dispatcher.d/20-chrony-dhcp
/usr/lib/NetworkManager/dispatcher.d/20-chrony-onoffline
/usr/lib/networkd-dispatcher
/usr/lib/networkd-dispatcher/off.d
/usr/lib/networkd-dispatcher/routable.d
/usr/lib/systemd
/usr/lib/systemd/ntp-units.d
/usr/lib/systemd/ntp-units.d/50-chrony.list
/usr/lib/systemd/scripts
/usr/lib/systemd/scripts/chronyd-starter.sh
/usr/lib/systemd/system
/usr/lib/systemd/system/chrony-dnssrv@.service
/usr/lib/systemd/system/chrony-dnssrv@.timer
/usr/lib/systemd/system/chrony-wait.service
/usr/lib/systemd/system/chrony.service
/usr/libexec
/usr/libexec/chrony
/usr/libexec/chrony/chrony-helper
/usr/sbin
/usr/sbin/chronyd
/usr/share
/usr/share/chrony
/usr/share/chrony/chrony.conf
/usr/share/chrony/chrony.keys
/usr/share/doc
/usr/share/doc/chrony
/usr/share/doc/chrony/FAQ.gz
/usr/share/doc/chrony/NEWS.Debian.gz
/usr/share/doc/chrony/NEWS.gz
/usr/share/doc/chrony/README.Debian
/usr/share/doc/chrony/README.container
/usr/share/doc/chrony/README.gz
/usr/share/doc/chrony/changelog.Debian.gz
/usr/share/doc/chrony/copyright
/usr/share/doc/chrony/examples
/usr/share/doc/chrony/examples/chrony.conf.example1
/usr/share/doc/chrony/examples/chrony.conf.example2
/usr/share/doc/chrony/examples/chrony.conf.example3
/usr/share/doc/chrony/examples/chronyd-restricted.service
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/chrony
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/chronyc.1.gz
/usr/share/man/man5
/usr/share/man/man5/chrony.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/chronyd.8.gz
/var
/var/lib
/var/lib/chrony
/usr/lib/networkd-dispatcher/off.d/chrony-onoffline
/usr/lib/networkd-dispatcher/routable.d/chrony-onoffline
