# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/iperf3.list
/.
/etc
/etc/init.d
/etc/init.d/iperf3
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/iperf3
/usr
/usr/bin
/usr/bin/iperf3
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/iperf3.service
/usr/share
/usr/share/doc
/usr/share/doc/iperf3
/usr/share/doc/iperf3/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/iperf3.1.gz
/usr/share/doc/iperf3/README.md.gz
/usr/share/doc/iperf3/changelog.Debian.gz
