# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/iperf3.conffiles
/etc/init.d/iperf3
/etc/ufw/applications.d/iperf3
