# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/ipp-usb.list
/.
/etc
/etc/ipp-usb
/etc/ipp-usb/ipp-usb.conf
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ipp-usb.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/71-ipp-usb.rules
/sbin/ipp-usb
/usr
/usr/share
/usr/share/doc
/usr/share/doc/ipp-usb
/usr/share/doc/ipp-usb/NEWS.Debian.gz
/usr/share/doc/ipp-usb/changelog.Debian.gz
/usr/share/doc/ipp-usb/copyright
/usr/share/ipp-usb
/usr/share/ipp-usb/quirks
/usr/share/ipp-usb/quirks/Canon.conf
/usr/share/ipp-usb/quirks/HP.conf
/usr/share/ipp-usb/quirks/README
/usr/share/ipp-usb/quirks/blacklist.conf
/usr/share/ipp-usb/quirks/default.conf
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ipp-usb.8.gz
