# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/modemmanager.list
/.
/etc
/etc/ModemManager
/etc/ModemManager/connection.d
/etc/ModemManager/fcc-unlock.d
/usr
/usr/bin
/usr/bin/mmcli
/usr/lib
/usr/lib/aarch64-linux-gnu
/usr/lib/aarch64-linux-gnu/ModemManager
/usr/lib/aarch64-linux-gnu/ModemManager/connection.d
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-altair-lte.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-anydata.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-broadmobi.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-cinterion.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-dell.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-dlink.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-ericsson-mbm.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-fibocom.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-foxconn.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-generic.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-gosuncn.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-haier.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-huawei.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-intel.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-iridium.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-linktop.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-longcheer.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-motorola.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-mtk-legacy.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-mtk.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-nokia-icera.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-nokia.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-novatel-lte.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-novatel.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-option-hso.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-option.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-pantech.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-qcom-soc.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-quectel.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-samsung.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-sierra-legacy.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-sierra.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-simtech.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-telit.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-thuraya.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-tplink.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-ublox.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-via.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-wavecom.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-x22x.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-plugin-zte.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-fibocom.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-foxconn.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-icera.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-novatel.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-option.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-sierra.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-telit.so
/usr/lib/aarch64-linux-gnu/ModemManager/libmm-shared-xmm.so
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/ModemManager.service
/usr/lib/udev
/usr/lib/udev/rules.d
/usr/lib/udev/rules.d/77-mm-broadmobi-port-types.rules
/usr/lib/udev/rules.d/77-mm-cinterion-port-types.rules
/usr/lib/udev/rules.d/77-mm-dell-port-types.rules
/usr/lib/udev/rules.d/77-mm-dlink-port-types.rules
/usr/lib/udev/rules.d/77-mm-ericsson-mbm.rules
/usr/lib/udev/rules.d/77-mm-fibocom-port-types.rules
/usr/lib/udev/rules.d/77-mm-foxconn-port-types.rules
/usr/lib/udev/rules.d/77-mm-gosuncn-port-types.rules
/usr/lib/udev/rules.d/77-mm-haier-port-types.rules
/usr/lib/udev/rules.d/77-mm-huawei-net-port-types.rules
/usr/lib/udev/rules.d/77-mm-linktop-port-types.rules
/usr/lib/udev/rules.d/77-mm-longcheer-port-types.rules
/usr/lib/udev/rules.d/77-mm-mtk-legacy-port-types.rules
/usr/lib/udev/rules.d/77-mm-mtk-port-types.rules
/usr/lib/udev/rules.d/77-mm-nokia-port-types.rules
/usr/lib/udev/rules.d/77-mm-qcom-soc.rules
/usr/lib/udev/rules.d/77-mm-qdl-device-blacklist.rules
/usr/lib/udev/rules.d/77-mm-quectel-port-types.rules
/usr/lib/udev/rules.d/77-mm-sierra.rules
/usr/lib/udev/rules.d/77-mm-simtech-port-types.rules
/usr/lib/udev/rules.d/77-mm-telit-port-types.rules
/usr/lib/udev/rules.d/77-mm-tplink-port-types.rules
/usr/lib/udev/rules.d/77-mm-ublox-port-types.rules
/usr/lib/udev/rules.d/77-mm-x22x-port-types.rules
/usr/lib/udev/rules.d/77-mm-zte-port-types.rules
/usr/lib/udev/rules.d/80-mm-candidate.rules
/usr/sbin
/usr/sbin/ModemManager
/usr/share
/usr/share/ModemManager
/usr/share/ModemManager/connection.available.d
/usr/share/ModemManager/connection.available.d/99-log-event
/usr/share/ModemManager/fcc-unlock.available.d
/usr/share/ModemManager/fcc-unlock.available.d/105b
/usr/share/ModemManager/fcc-unlock.available.d/1199
/usr/share/ModemManager/fcc-unlock.available.d/14c3
/usr/share/ModemManager/fcc-unlock.available.d/1eac
/usr/share/ModemManager/fcc-unlock.available.d/2c7c
/usr/share/ModemManager/mm-foxconn-t77w968-carrier-mapping.conf
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/mmcli
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.ModemManager1.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.ModemManager1.conf
/usr/share/doc
/usr/share/doc/modemmanager
/usr/share/doc/modemmanager/AUTHORS
/usr/share/doc/modemmanager/NEWS.Debian.gz
/usr/share/doc/modemmanager/NEWS.gz
/usr/share/doc/modemmanager/README.Debian
/usr/share/doc/modemmanager/README.md
/usr/share/doc/modemmanager/copyright
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/22x22
/usr/share/icons/hicolor/22x22/apps
/usr/share/icons/hicolor/22x22/apps/ModemManager.png
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/modemmanager
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mmcli.1.gz
/usr/share/man/man8
/usr/share/man/man8/ModemManager.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.ModemManager1.policy
/usr/share/ModemManager/fcc-unlock.available.d/03f0:4e1d
/usr/share/ModemManager/fcc-unlock.available.d/105b:e0ab
/usr/share/ModemManager/fcc-unlock.available.d/105b:e0c3
/usr/share/ModemManager/fcc-unlock.available.d/1199:9079
/usr/share/ModemManager/fcc-unlock.available.d/14c3:4d75
/usr/share/ModemManager/fcc-unlock.available.d/1eac:1001
/usr/share/ModemManager/fcc-unlock.available.d/2c7c:030a
/usr/share/ModemManager/fcc-unlock.available.d/413c:81a3
/usr/share/ModemManager/fcc-unlock.available.d/413c:81a8
/usr/share/doc/modemmanager/changelog.Debian.gz
