# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/ntpsec.list
/etc
/etc/apparmor.d
/etc/apparmor.d/tunables
/etc/apparmor.d/tunables/ntpd
/etc/apparmor.d/usr.sbin.ntpd
/etc/cron.d
/etc/cron.d/ntpsec
/etc/default
/etc/default/ntpsec
/etc/dhcp
/etc/dhcp/dhclient-exit-hooks.d
/etc/dhcp/dhclient-exit-hooks.d/ntpsec
/etc/init.d
/etc/init.d/ntpsec
/etc/letsencrypt
/etc/letsencrypt/renewal-hooks
/etc/letsencrypt/renewal-hooks/deploy
/etc/letsencrypt/renewal-hooks/deploy/ntpsec
/etc/ntpsec
/etc/ntpsec/ntp.conf
