# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/sane-utils.conffiles
/etc/default/saned
/etc/init.d/saned
/etc/logrotate.d/sane-utils
/etc/sane.d/saned.conf
