# houzzkit-f1-rootfs-u2404/ubuntu-rootfs/var/lib/dpkg/info/usbmuxd.list
/.
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/usbmuxd.service
/usr/lib/udev
/usr/lib/udev/rules.d
/usr/lib/udev/rules.d/39-usbmuxd.rules
/usr/sbin
/usr/sbin/usbmuxd
/usr/share
/usr/share/doc
/usr/share/doc/usbmuxd
/usr/share/doc/usbmuxd/AUTHORS
/usr/share/doc/usbmuxd/README.md.gz
/usr/share/doc/usbmuxd/changelog.Debian.gz
/usr/share/doc/usbmuxd/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/usbmuxd.8.gz
/usr/share/metainfo
/usr/share/metainfo/org.libimobiledevice.usbmuxd.metainfo.xml
