# houzzkit-f1-bsp-k6199/yocto/meta-rockchip/recipes-bsp/rockchip-npu/files/rknn_server.service
[Unit]
Description=start rknn_server service

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=rknn_server

[Install]
WantedBy=local-fs.target
