# houzzkit-f1-bsp-k6199/yocto/meta-rockchip/overlay-firefly/lib/systemd/system/ec20.service
[Unit]
Description=Quectel networking service
Before=network.target
Wants=network.target
BindsTo=dev-ttyUSB1.device
After=dev-ttyUSB1.device

[Service]
Type=simple
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/quectel-CM -f /tmp/quectel-CM_log
Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
